# udev-rules/81-brightness.rules
# Restore the brightness on startup:
ACTION=="add", SUBSYSTEM=="backlight", KERNEL=="intel_backlight", RUN+="/etc/udev/brightness.sh restore"